# server/upstart.conf
start on filesystem and started networking
respawn
console log
chdir /home/poscar/webdevelopment.rocks/server
exec /usr/bin/nodejs index.js
